# project/NW5002/files/etc/hotplug.d/block/20-mmc
#!/bin/sh

mkdir -p /tmp/mnt/SD-disk-p4
#mkfs.vfat /dev/mmcblk0p4
mount -t vfat -o fmask=000,dmask=000,iocharset=utf8 /dev/mmcblk0p4 /tmp/mnt/SD-disk-p4
