# project/nw5025F_ac_xinao/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

check_reset &
/usr/sbin/check_shair.sh &
updatefw &
erase_boot_flag auto
#init_network.sh
#route add -net 224.0.0.0 netmask 255.255.255.0 wl0.1
#airplay -a airdisk -d
exit 0
