# project/orico_round/package/dm_discovery/files/dm_discovery.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009-2012 OpenWrt.org

START=90

RUN_D=/var/run/msg_server

start() {
        dm_discovery &
}
stop() {
        killall dm_discovery
}
restart() {
		killall dm_discovery
		dm_discovery &
}
