# project/airphoto/files/etc/hotplug.d/block/20-mmc
#!/bin/sh

mkdir -p /tmp/mnt/SD-disk-p4

mount -t ext4 /dev/mmcblk0p4 /tmp/mnt/SD-disk-p4
